# ingest-worker/deploy/watchlater-ingest.service
[Unit]
Description=WatchLater ingest worker (transcript + lesson generation)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=watchlater
Group=watchlater
WorkingDirectory=/opt/watchlater/ingest-worker
EnvironmentFile=/etc/watchlater/ingest.env
Environment=HOME=/var/lib/watchlater
StateDirectory=watchlater
ExecStart=/usr/local/bin/bun run src/index.ts
Restart=always
RestartSec=3

# Hardening — the worker only makes outbound HTTPS calls and listens on
# localhost; it needs no privileges beyond reading its working directory.
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
